# ddf-common-script/src/main/resources/shell/supervistor.ini
[program:xiuxian]
directory=/opt/services/xiu-xian-config-server
command=sh /opt/services/xiu-xian-config-server/simple_deploy.sh 8080
autostart=true
autorestart=true
priority=5
user=root
stdout_logfile=/dev/null
redirect_stderr=true
loglevel=info
stopasgroup=true
killasgroup=true
